# docker/engine/sync-vendor/myprogram.conf
[program:myprogram]
command=/usr/local/bin/sync_script.sh
autorestart=true
startsecs=5
startretries=3
